# config/init-uwsgi.conf
# simple uWSGI script
description "Django-uWSGI"
start on runlevel [2345]
stop on runlevel [06]

exec uwsgi --ini /home/ubuntu/pythia-server/config/django-uwsgi.ini
